# macos/com.assistant.hotkey.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.assistant.hotkey</string>
    <key>ProgramArguments</key>
    <array>
        <string>__ASSISTANT_DIR__/dist/assistant-hotkey</string>
    </array>
    <key>EnvironmentVariables</key>
    <dict>
        <key>ASSISTANT_DIR</key>
        <string>__ASSISTANT_DIR__</string>
    </dict>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/tmp/assistant-hotkey.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/assistant-hotkey.log</string>
</dict>
</plist>
